# fs/etc/init.d/set-timezone
#!/bin/bash

### BEGIN INIT INFO
# Default-Start:	2
# Default-Stop:		
### END INIT INFO

timedatectl set-timezone Asia/Tehran

update-rc.d set-timezone remove
