# fs/etc/init.d/update-initramfs
#!/bin/bash

### BEGIN INIT INFO
# Default-Start:	2
# Default-Stop:		
### END INIT INFO

set -e

update-initramfs -u

update-rc.d update-initramfs remove
